# systemd/hr-sync.service
[Unit]
Description=HR EMS — CSV → PostgreSQL 동기화
After=network.target hr-postgres.service
Requires=hr-postgres.service

[Service]
Type=simple
User=dt
WorkingDirectory=/home/dt/hr_1.1
EnvironmentFile=/home/dt/hr_1.1/site.env
ExecStart=/home/dt/hr_1.1/venv/bin/python3 bin/sync.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
